# services/emborg-backup.timer
[Unit]
Description=Run emborg backup @@DATE@@

[Timer]
OnCalendar=@@DATE@@
Persistent=true

[Install]
WantedBy=timers.target